# chef/templates/default/systemd.service.erb
[Unit]
Description=<%= @name %>
<% if @options['after_target'] && ! @options['after_target'].empty? -%>
After=<%= @options['after_target'] %><%= '.target' unless @options['after_target'].include?('.') %>
<% end -%>

[Service]
Environment=<%= @environment.map {|key, val| %Q{"#{key}=#{val}"} }.join(' ') %>
ExecStart=<%= @command %>
ExecReload=/bin/kill -<%= @reload_signal %> $MAINPID
KillSignal=<%= @stop_signal %>
User=<%= @user %>
WorkingDirectory=<%= @directory %>
Restart=<%= @restart_mode %>

[Install]
WantedBy=multi-user.target
